# systemd/nfts-bot.service
# systemd/nfts-bot.service
[Unit]
Description=NFT Collection Stats Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nft-discord-bots
Environment=PYTHONPATH=/home/ubuntu/nft-discord-bots
ExecStart=/home/ubuntu/nft-discord-bots/venv/bin/python -c 'import os; from bots.nfts_bot import NFTsBot; NFTsBot().run(os.getenv("NFTS_BOT_TOKEN"))'
Restart=always

[Install]
WantedBy=multi-unit.target