# systemd/offers-bot.service
# systemd/offers-bot.service
[Unit]
Description=NFT Offers Stats Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nft-discord-bots
Environment=PYTHONPATH=/home/ubuntu/nft-discord-bots
ExecStart=/home/ubuntu/nft-discord-bots/venv/bin/python -c 'import os; from bots.offers_bot import OffersBot; OffersBot().run(os.getenv("OFFERS_BOT_TOKEN"))'
Restart=always

[Install]
WantedBy=multi-unit.target