# systemd/users-bot.service
# systemd/users-bot.service
[Unit]
Description=NFT Users Stats Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nft-discord-bots
Environment=PYTHONPATH=/home/ubuntu/nft-discord-bots
ExecStart=/home/ubuntu/nft-discord-bots/venv/bin/python -c 'import os; from bots.users_bot import UsersBot; UsersBot().run(os.getenv("USERS_BOT_TOKEN"))'
Restart=always

[Install]
WantedBy=multi-unit.target