# systemd/users-bot.service
# systemd/users-bot.service
[Unit]
Description=NFT Users Stats Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/discord-bot/bots
Environment=PYTHONPATH=/home/ubuntu/discord-bot/bots
ExecStart=/home/ubuntu/discord-bot/venv/bin/python -c 'from users_bot import UsersBot; UsersBot().run(os.getenv("USERS_BOT_TOKEN"))'
Restart=always

[Install]
WantedBy=multi-unit.target